# copilot-api-proxy.socket
[Unit]
Description=GitHub Copilot API Proxy Socket v2.0
Documentation=https://github.com/dcx7c5/copilot-api-proxy
# NOTE: This socket unit is only used when SERVER__MODE=unix_socket.
#       In tls_socket mode (default) the service binds directly to TCP 8443
#       and this unit must remain DISABLED.
#
#       To switch to unix_socket mode:
#         1. Set SERVER__MODE=unix_socket in /etc/copilot-api-proxy/config.env
#         2. sudo systemctl enable --now copilot-api-proxy.socket
#         3. sudo systemctl disable copilot-api-proxy.service  (service started by socket)

[Socket]
# Unix-domain socket for unix_socket mode
ListenStream=/run/copilot-api-proxy/copilot-proxy.sock
SocketUser=copilot-api-proxy
SocketGroup=copilot-api-proxy
SocketMode=0660
DirectoryMode=0755

# Socket tuning
Backlog=256
KeepAlive=yes
KeepAliveTimeSec=30
KeepAliveIntervalSec=5
KeepAliveProbes=3
MaxConnections=100
MaxConnectionsPerSource=10
ReceiveBuffer=65535
SendBuffer=65535

# Avoid hammering a crashing service
TriggerLimitIntervalSec=30s
TriggerLimitBurst=5

# Cleanup
RemoveOnStop=yes
FreeBind=no

[Install]
# Do NOT add WantedBy=sockets.target here — enable manually only when
# unix_socket mode is configured.
